# roles/utils/updateable/templates/anacrontab
# /etc/anacrontab: configuration file for anacron

# See anacron(8) and anacrontab(5) for details.

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
HOME={{ user.homedir }}
LOGNAME={{ user.name}}

# These replace cron's entries
1	3	{{ user.name }}.cron.daily	shall_do_update && run-parts --report {{ update_daily }}
7	8	{{ user.name }}.cron.weekly	shall_do_update && run-parts --report {{ update_weekly }}
@monthly	8	{{ user.name }}.cron.monthly	shall_do_update && run-parts --report {{ update_monthly }}
